# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/etc/update-motd.d/00-header
#!/bin/sh
#

[ -r /etc/lsb-release ] && . /etc/lsb-release

#printf "\n"
echo "\e[92m$(cat /usr/share/doc/femtofox/logo)\033[0m"

if [ -z "$DISTRIB_DESCRIPTION" ] && [ -x /usr/bin/lsb_release ]; then
	# Fall back to using the very slow lsb_release utility
	DISTRIB_DESCRIPTION=$(lsb_release -s -d)
fi

[ -r /etc/lsb-release ] && . /etc/lsb-release
foxbuntu_version=$(femto-utils.sh -v)
spaces=$((21 - ${#foxbuntu_version}))

printf "                     For settings, run \`\033[0;31msudo femto-config\033[0m\`\n"
printf "                  Made for Femtofox  https://www.femtofox.com\n"
printf "               $foxbuntu_version%*s%s(%s)\n\n" $spaces "" "$DISTRIB_DESCRIPTION" "$(uname -r)"

printf "IP address: $(hostname -I | sed 's/ / /g') Uptime:$(uptime -p | awk '{$1=""; print $0}' | sed -e 's/ day\b/d/g' -e 's/ hour\b/h/g' -e 's/ hours\b/h/g' -e 's/ minute\b/m/g' -e 's/ minutes\b/m/g' | sed 's/,//g')\n"
